# files/etc/crontabs/root
0 3 * * 0 /etc/init.d/dnscrypt-proxy update_servers >/dev/null 2>&1
