# files/etc/uci-defaults/60-stunnel-migration
#!/bin/sh
# One-shot migration for stunnel+tun2socks architecture.
# Removes obsolete openconnect UCI keys and initializes new ones.

uci -q delete vipin.vpn.enabled
uci -q delete vipin.vpn.cert_pin

# Seed auth_status=ok so first boot does not show an error banner
# before vipin-auth-check has a chance to run.
[ -z "$(uci -q get vipin.vpn.auth_status)" ] && uci set vipin.vpn.auth_status='ok'

# DO NOT seed split_mode here. Earlier versions defaulted to 'forward'
# which preempted vipin-auto-login's country-aware default (overseas →
# reverse, CN → forward) — this script runs at procd preinit, auto-login
# fires later on WAN ifup, and auto-login only writes when UCI is empty.
# Leaving split_mode unset lets auto-login pick, and if the user never
# runs auto-login (no credentials baked) LuCI still surfaces all three
# radio options and prompts for a choice.

uci commit vipin

# Also remove the legacy openconnect cert pin file if present.
rm -f /etc/vipin/cert_pin

# Stop and disable OpenWrt's default system stunnel. The stock package
# ships with a dummy loopback [dummy] section (localhost:6000→6001) that
# nobody uses but still costs ~7 MB RSS on a 128 MB router. Our VPN
# stack runs its own stunnel via /etc/init.d/vipin-vpn, so disable the
# system one to free memory.
/etc/init.d/stunnel stop 2>/dev/null
/etc/init.d/stunnel disable 2>/dev/null

# Same deal with system stubby: package auto-starts its own daemon on
# 127.0.0.1:53 bound to cloudflare/quad9 DoT. Our vipin-vpn manages a
# different stubby instance via procd when overseas-to-CN DNS is needed.
/etc/init.d/stubby stop 2>/dev/null
/etc/init.d/stubby disable 2>/dev/null

exit 0
