# files/etc/uci-defaults/luci-dns-dnscrypt
#!/bin/sh
# Enable dnscrypt-proxy service at boot — do NOT set DNS policy here.
#
# Rationale: vipin-vpn's configure_dns() / teardown_dns() is the single
# authority for dnsmasq upstream selection. Pre-VPN the policy must be
# WAN DHCP-learned DNS (otherwise VPN can't resolve its server domain
# and dnscrypt can't reach DoH upstreams — chicken-and-egg deadlock).
# zzz-defaults leaves dnsmasq with no explicit server so it reads
# resolv.conf.auto. vipin-vpn flips to 127.0.0.1#5353 when VPN is up.

/etc/init.d/dnscrypt-proxy enable 2>/dev/null
/etc/init.d/dnscrypt-proxy start 2>/dev/null

rm -f /etc/uci-defaults/luci-dns-dnscrypt
exit 0
