# files/etc/uci-defaults/luci-dns-dnscrypt
#!/bin/sh
# Configure dnsmasq to use dnscrypt-proxy as upstream DNS

dnsmasq_cfg="/etc/config/dhcp"

# Wait for base system to initialize
sleep 2

# Ensure dnsmasq config exists
[ ! -f "$dnsmasq_cfg" ] && {
    cat > "$dnsmasq_cfg" << 'EOF'
config dnsmasq
	option domainneeded '1'
	option localise_queries '1'
	option rebind_protection '1'
	option rebind_localhost '1'
	option local '/lan/'
	option domain 'lan'
	option expandhosts '1'
	option authoritative '1'
	option readethers '1'
	option leasefile '/tmp/dhcp.leases'
	option resolvfile '/tmp/resolv.conf.d/resolv.conf.auto'
	option localservice '1'

config dhcp 'lan'
	option interface 'lan'
	option start '100'
	option limit '150'
	option leasetime '12h'

config dhcp 'wan'
	option interface 'wan'
	option ignore '1'
EOF
}

# Configure dnsmasq to use dnscrypt-proxy
uci -q del_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'
uci -q add_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'
uci -q commit dhcp

/etc/init.d/dnsmasq restart 2>/dev/null

# Enable dnscrypt-proxy on boot
/etc/init.d/dnscrypt-proxy enable 2>/dev/null
/etc/init.d/dnscrypt-proxy start 2>/dev/null

rm -f /etc/uci-defaults/luci-dns-dnscrypt
exit 0
