# files/etc/uci-defaults/luci-i18n-vipin
#!/bin/sh
# Install custom i18n locale files (.po text files) for LuCI VPN/DNSCrypt modules
# OpenWRT LuCI reads .po files directly if .lmo is not available

DEST_DIR="/usr/share/luci/i18n"
mkdir -p "$DEST_DIR"

for locale_dir in /usr/share/luci/locale/*/; do
    lang=$(basename "$locale_dir")
    for po_file in "$locale_dir"/LC_MESSAGES/*.po 2>/dev/null; do
        [ -f "$po_file" ] || continue
        module=$(basename "$po_file" .po)
        install -m 0644 "$po_file" "$DEST_DIR/${module}.${lang}.po"
    done
done

rm -f /etc/uci-defaults/luci-i18n-vipin

exit 0
