# files/etc/uci-defaults/zzz-defaults
#!/bin/sh

# Set hostname
uci set system.@system[0].hostname='202604'
uci commit system

# Set root password (reads from UCI, fallback to site domain)
SITE_DOMAIN=$(uci get vipin.vpn.site_url 2>/dev/null | sed 's|https\?://||;s|/.*||' || echo "www.anyfq.com")
echo -e "${SITE_DOMAIN}\n${SITE_DOMAIN}" | passwd root

# Set LAN IP
uci set network.lan.ipaddr='192.168.11.1'
uci set network.lan.netmask='255.255.255.0'
uci commit network

# Force generate wireless config if not exists
[ ! -f /etc/config/wireless ] && wifi detect > /etc/config/wireless

# Wait for wireless config to be generated
sleep 2

# Enable all radios
for radio in $(uci show wireless | grep "=wifi-device" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${radio}.disabled='0'
done

# Configure all wifi-iface
for iface in $(uci show wireless | grep "=wifi-iface" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${iface}.ssid="${SITE_DOMAIN}"
    uci set wireless.${iface}.encryption='psk2'
    uci set wireless.${iface}.key="${SITE_DOMAIN}"
    uci set wireless.${iface}.disabled='0'
done

# Commit and bring up WiFi
uci commit wireless
wifi up

# Reload network
/etc/init.d/network reload

exit 0
