# files/etc/uci-defaults/zzz-defaults
#!/bin/sh

# Set hostname
uci set system.@system[0].hostname='202604'
uci commit system

# Set root password (reads from UCI, fallback to site domain)
SITE_DOMAIN=$(uci get vipin.vpn.site_url 2>/dev/null | sed 's|https\?://||;s|/.*||' || echo "www.anyfq.com")
echo -e "${SITE_DOMAIN}\n${SITE_DOMAIN}" | passwd root

# Set LAN IP
uci set network.lan.ipaddr='192.168.11.1'
uci set network.lan.netmask='255.255.255.0'
# Never accept DNS from WAN DHCP (LAN must not resolve via ISP DNS)
uci set network.wan.peerdns='0'
uci -q delete network.wan.dns
uci set network.wan6.peerdns='0'
uci -q delete network.wan6.dns
uci commit network

# Disable LAN-side IPv6 entirely — prevents IPv6 bypassing VPN split-tunnel.
# Geo-locked services (e.g. Tencent Video) would otherwise see real WAN IPv6.
uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.dhcpv6='disabled'
uci -q delete dhcp.lan.ra_slaac
uci -q delete dhcp.lan.ra_flags

# DNS policy: dnsmasq never reads resolv.conf.auto (no WAN DNS fallback).
# Default upstream = dnscrypt-proxy on 127.0.0.1#5353.
# vipin-vpnc-script swaps upstream to VPN-pushed DNS on connect.
uci set dhcp.@dnsmasq[0].noresolv='1'
uci set dhcp.@dnsmasq[0].resolvfile='/dev/null'
# Belt-and-braces: even with LAN IPv6 disabled, don't serve AAAA to clients
# (a device that enables IPv6 manually still shouldn't get geo-leak AAAA).
uci set dhcp.@dnsmasq[0].filter_aaaa='1'
uci -q delete dhcp.@dnsmasq[0].server
uci add_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'
uci commit dhcp

# Force generate wireless config if not exists
[ ! -f /etc/config/wireless ] && wifi detect > /etc/config/wireless

# Wait for wireless config to be generated
sleep 2

# Enable all radios
for radio in $(uci show wireless | grep "=wifi-device" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${radio}.disabled='0'
done

# Configure all wifi-iface
for iface in $(uci show wireless | grep "=wifi-iface" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${iface}.ssid="${SITE_DOMAIN}"
    uci set wireless.${iface}.encryption='psk2'
    uci set wireless.${iface}.key="${SITE_DOMAIN}"
    uci set wireless.${iface}.disabled='0'
done

# Commit and bring up WiFi
uci commit wireless
wifi up

# Reload network
/etc/init.d/network reload

exit 0
