# files/etc/uci-defaults/zzz-defaults
#!/bin/sh

# Set hostname
uci set system.@system[0].hostname='202604'
uci commit system

# Set root password (reads from UCI, fallback to site domain)
SITE_DOMAIN=$(uci get vipin.vpn.site_url 2>/dev/null | sed 's|https\?://||;s|/.*||' || echo "www.anyfq.com")
echo -e "${SITE_DOMAIN}\n${SITE_DOMAIN}" | passwd root

# Set LAN IP
uci set network.lan.ipaddr='192.168.11.1'
uci set network.lan.netmask='255.255.255.0'
# peerdns=1 so netifd writes ISP-provided DNS into /tmp/resolv.conf*.auto.
# dnsmasq still stays noresolv=1 (configure_dns is single source of truth
# for the LAN upstream); this just makes ISP DNS *available* for scenarios
# that need it (VPN-down fallback, overseas-router + non-CN-server).
uci set network.wan.peerdns='1'
uci -q delete network.wan.dns
uci commit network

# Disable IPv6 across the whole router — no kernel v6, no WAN6, no LAN6,
# no odhcpd. Geo-locked services would otherwise see the real WAN IPv6
# regardless of split-tunnel IPv4 rules.
uci -q delete network.wan6
uci -q delete network.globals.ula_prefix
uci set network.lan.ipv6='0'
uci set network.lan.delegate='0'
uci commit network

uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.dhcpv6='disabled'
uci set dhcp.lan.ndp='disabled'
uci -q delete dhcp.lan.ra_slaac
uci -q delete dhcp.lan.ra_flags
uci set dhcp.odhcpd.maindhcp='0'

# DNS policy at first-boot: dnscrypt-proxy is the preferred upstream, but
# dnsmasq IS allowed to read resolv.conf.auto (WAN DHCP-learned DNS) as
# fallback. This breaks the pre-VPN chicken-and-egg deadlock seen on
# first-boot when the clock is wrong: dnscrypt TLS fails → no DNS → no
# NTP → clock stays wrong. With a WAN-DNS fallback, name resolution
# works even before dnscrypt/NTP recover.
#
# VPN is not yet connected at this point, so there is no CN geo-spoof
# leak to worry about — that only matters once tun0 is up and CN video
# CDNs are being hit. vipin-vpn start_service's configure_dns() then
# rewrites the server list + sets noresolv=1 to enforce the final
# policy (CN DNS for overseas→CN-server, dnscrypt for CN-router, etc.).
uci -q delete dhcp.@dnsmasq[0].noresolv
uci -q delete dhcp.@dnsmasq[0].resolvfile
# Belt-and-braces: even with LAN IPv6 disabled, don't serve AAAA to clients
# (a device that enables IPv6 manually still shouldn't get geo-leak AAAA).
uci set dhcp.@dnsmasq[0].filter_aaaa='1'
# 10k cache floors first-query-fail impact: with only 150 (OpenWrt default)
# any DNS-miss burst starves the cache before repeat queries can land.
# 10k entries ≈ 8 MB extra dnsmasq RSS — acceptable on 128 MB routers.
uci set dhcp.@dnsmasq[0].cachesize='10000'
# Silence per-DHCPACK/REQUEST log lines. Misbehaving clients (Android TV
# boxes with residential-proxy SDKs) can DHCPREQUEST every 3s; with
# default logging that's 1200 syslog writes/hour per client and
# perceptible CPU load on MT7628. DHCP-level errors still log.
uci set dhcp.@dnsmasq[0].quietdhcp='1'
uci -q delete dhcp.@dnsmasq[0].server
uci add_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'
uci commit dhcp

# Kernel-level IPv6 disable — nothing on this box will have an IPv6 address
cat > /etc/sysctl.d/99-disable-ipv6.conf << 'SYSCTL'
net.ipv6.conf.all.disable_ipv6=1
net.ipv6.conf.default.disable_ipv6=1
net.ipv6.conf.lo.disable_ipv6=1
SYSCTL
sysctl -p /etc/sysctl.d/99-disable-ipv6.conf >/dev/null 2>&1

# Stop/disable odhcpd — no RA, no DHCPv6, nothing
/etc/init.d/odhcpd disable 2>/dev/null
/etc/init.d/odhcpd stop 2>/dev/null

# Force generate wireless config if not exists
[ ! -f /etc/config/wireless ] && wifi detect > /etc/config/wireless

# Wait for wireless config to be generated
sleep 2

# Enable all radios
for radio in $(uci show wireless | grep "=wifi-device" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${radio}.disabled='0'
done

# Configure all wifi-iface
for iface in $(uci show wireless | grep "=wifi-iface" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${iface}.ssid="${SITE_DOMAIN}"
    uci set wireless.${iface}.encryption='psk2'
    uci set wireless.${iface}.key="${SITE_DOMAIN}"
    uci set wireless.${iface}.disabled='0'
done

# Commit and bring up WiFi
uci commit wireless
wifi up

# Reload network
/etc/init.d/network reload

exit 0
