# files/etc/uci-defaults/zzz-defaults
#!/bin/sh

# Set hostname
uci set system.@system[0].hostname='202604'
uci commit system

# Set root password
echo -e 'www.anyfq.com\nwww.anyfq.com' | passwd root

# Set LAN IP
uci set network.lan.ipaddr='192.168.11.1'
uci set network.lan.netmask='255.255.255.0'
uci commit network

# Force generate wireless config if not exists
[ ! -f /etc/config/wireless ] && wifi detect > /etc/config/wireless

# Wait for wireless config to be generated
sleep 2

# Enable all radios and set WiFi defaults
for radio in $(uci show wireless | grep "=wifi-device" | cut -d. -f2); do
    # Enable radio
    uci set wireless.${radio}.disabled='0'
done

# Configure all wifi-iface
for iface in $(uci show wireless | grep "=wifi-iface" | cut -d. -f2); do
    uci set wireless.${iface}.ssid='www.anyfq.com'
    uci set wireless.${iface}.encryption='psk2'
    uci set wireless.${iface}.key='www.anyfq.com'
    uci set wireless.${iface}.disabled='0'
done

# Commit and bring up WiFi
uci commit wireless
wifi up

# Reload network
/etc/init.d/network reload

exit 0
