# files/etc/uci-defaults/zzz-defaults
#!/bin/sh

# dnsmasq ujail-mount fix: /etc/dnsmasq.conf ships a `conf-dir=/etc/dnsmasq.d/`
# directive (OpenWrt default), and procd's ujail tries to bind-mount that
# directory into the dnsmasq jail. If the directory is missing at first boot
# (sysupgrade or fresh flash with no overlay yet), the mount fails and dnsmasq
# fatals with "cannot access directory /etc/dnsmasq.d/" → procd crash loop →
# DHCP/DNS dark. Pre-create the directory so the mount always succeeds.
mkdir -p /etc/dnsmasq.d

# Set hostname
uci set system.@system[0].hostname='202604'
uci commit system

# Set root password (reads from UCI, fallback to site domain)
SITE_DOMAIN=$(uci get vipin.vpn.site_url 2>/dev/null | sed 's|https\?://||;s|/.*||' || echo "www.anyfq.com")
echo -e "${SITE_DOMAIN}\n${SITE_DOMAIN}" | passwd root

# Set LAN IP
uci set network.lan.ipaddr='192.168.11.1'
uci set network.lan.netmask='255.255.255.0'
# peerdns=1 so netifd writes ISP-provided DNS into /tmp/resolv.conf*.auto.
# dnsmasq still stays noresolv=1 (configure_dns is single source of truth
# for the LAN upstream); this just makes ISP DNS *available* for scenarios
# that need it (VPN-down fallback, overseas-router + non-CN-server).
uci set network.wan.peerdns='1'
uci -q delete network.wan.dns
uci commit network

# Disable IPv6 across the whole router — no kernel v6, no WAN6, no LAN6,
# no odhcpd. Geo-locked services would otherwise see the real WAN IPv6
# regardless of split-tunnel IPv4 rules.
uci -q delete network.wan6
uci -q delete network.globals.ula_prefix
uci set network.lan.ipv6='0'
uci set network.lan.delegate='0'
uci commit network

uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.dhcpv6='disabled'
uci set dhcp.lan.ndp='disabled'
uci -q delete dhcp.lan.ra_slaac
uci -q delete dhcp.lan.ra_flags
uci set dhcp.odhcpd.maindhcp='0'

# DNS policy at first-boot: WAN DHCP-learned DNS only. vipin-vpn is the
# single authority that switches dnsmasq upstream (configure_dns points
# at un.fanq.in:3333 over the hysteria2 tunnel once VPN is up; teardown
# falls back to a pinned WAN DNS).
#
# Pre-VPN: WAN DHCP DNS resolves vipin-vpn's bootstrap targets. Once
# hysteria2 tunnel is established, dnsmasq forwards to un.fanq.in:3333
# (socat → un.fanq.in's local dnsmasq) which terminates on un's
# in-tunnel DNS resolver, avoiding GFW DNS interception entirely.
#
# With noresolv unset + no explicit server here, dnsmasq reads
# /tmp/resolv.conf.d/resolv.conf.auto (populated by WAN DHCP because
# peerdns=1 is set above). Name resolution works from the moment the
# WAN link is up, regardless of VPN state.
uci -q delete dhcp.@dnsmasq[0].noresolv
uci -q delete dhcp.@dnsmasq[0].resolvfile
# Belt-and-braces: even with LAN IPv6 disabled, don't serve AAAA to clients
# (a device that enables IPv6 manually still shouldn't get geo-leak AAAA).
uci set dhcp.@dnsmasq[0].filter_aaaa='1'
# 10k cache floors first-query-fail impact: with only 150 (OpenWrt default)
# any DNS-miss burst starves the cache before repeat queries can land.
# 10k entries ≈ 8 MB extra dnsmasq RSS — acceptable on 128 MB routers.
uci set dhcp.@dnsmasq[0].cachesize='10000'
# Silence per-DHCPACK/REQUEST log lines. Misbehaving clients (Android TV
# boxes with residential-proxy SDKs) can DHCPREQUEST every 3s; with
# default logging that's 1200 syslog writes/hour per client and
# perceptible CPU load on MT7628. DHCP-level errors still log.
uci set dhcp.@dnsmasq[0].quietdhcp='1'
# DO NOT set dhcp.@dnsmasq[0].server at first-boot — that would override
# WAN DHCP DNS before vipin-vpn can swap to un.fanq.in:3333 (which itself
# requires the tunnel to be up to be reachable). vipin-vpn owns the
# server list at runtime.
uci -q delete dhcp.@dnsmasq[0].server
uci commit dhcp

# Kernel-level IPv6 disable — nothing on this box will have an IPv6 address
cat > /etc/sysctl.d/99-disable-ipv6.conf << 'SYSCTL'
net.ipv6.conf.all.disable_ipv6=1
net.ipv6.conf.default.disable_ipv6=1
net.ipv6.conf.lo.disable_ipv6=1
SYSCTL
sysctl -p /etc/sysctl.d/99-disable-ipv6.conf >/dev/null 2>&1

# Stop/disable odhcpd — no RA, no DHCPv6, nothing
/etc/init.d/odhcpd disable 2>/dev/null
/etc/init.d/odhcpd stop 2>/dev/null

# Force generate wireless config if not exists
[ ! -f /etc/config/wireless ] && wifi detect > /etc/config/wireless

# Wait for wireless config to be generated
sleep 2

# Enable all radios + per-band tuning.
#  noscan=1  (2.4G only) skip HT40 coexistence scan — weak SoCs (MT7628) would
#            otherwise fall back to HT20 whenever any neighbor AP is seen.
#  htmode=HT40 (2.4G only) force 40 MHz channel width for 2× throughput.
#            Without this, wifi-detect on some chipsets defaults to HT20.
#            On 5G we leave htmode untouched (VHT80 / HE80 defaults are fine).
#            noscan + HT40 pair cleanly: 40 MHz stays sticky even with noisy
#            neighbors (our site is single-AP on ch1, DDXS is on ch11 — no
#            real coexistence to yield to).
for radio in $(uci show wireless | grep "=wifi-device" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${radio}.disabled='0'
    band=$(uci -q get wireless.${radio}.band)
    if [ "$band" = "2g" ]; then
        uci set wireless.${radio}.noscan='1'
        uci set wireless.${radio}.htmode='HT40'
    fi
done

# Configure all wifi-iface.
# Stability tweaks learned the hard way on MT7628 (GL-MT300N-V2) @ ch1:
# Android TVs + weak 2.4G radios cycled every ~15 min (beacon-loss roam),
# and channel busy time was dominated by 1 Mbps mgmt/broadcast frames.
# Applied universally — all safe no-op on healthy setups.
#  disassoc_low_ack=0      don't kick STA on low ACK (false-positive on weak SNR)
#  max_inactivity=600      raise hostapd idle timeout from 300s default
#  skip_inactivity_poll=1  don't send null-data poll (triggers re-auth on some TVs)
#  wpa_group_rekey=86400   GTK rekey 10min→24h (reduces handshake churn)
#  basic_rate 6/12/24 Mbps drop 11b legacy (1/2/5.5/11). Frees ~15% airtime on
#                          2.4G from management frames; no-op on 5G (no 11b).
#  mcast_rate=24000        broadcasts/multicasts at 24 Mbps instead of 1 Mbps.
#                          Every beacon, DHCP ACK, ARP reply saves 24× airtime.
for iface in $(uci show wireless | grep "=wifi-iface" | sed 's/.*\.//' | sed 's/=.*//'); do
    uci set wireless.${iface}.ssid="${SITE_DOMAIN}"
    # WPA2+WPA3 mixed: modern STAs (iPhone 6S+/iOS 13+, Android 10+) negotiate
    # SAE → forward-secrecy + offline-brute-force immune handshake; old STAs
    # (TV STBs, iPhone 4/5) fall back to WPA2-PSK CCMP. Same passphrase reused
    # via sae_password (hostapd derives from wpa_passphrase by default).
    uci set wireless.${iface}.encryption='sae-mixed'
    uci set wireless.${iface}.key="${SITE_DOMAIN}"
    uci set wireless.${iface}.disabled='0'
    uci set wireless.${iface}.disassoc_low_ack='0'
    # No-kick suite for fragile STBs (verified: TV box with broken
    # wpa_supplicant auto-reconnect was getting deauthed ~19×/30h on a
    # 1800s inactivity setting; user had to manually toggle WiFi every time).
    # max_inactivity 31536000 (1y) effectively disables hostapd's STA
    # idle-eviction; wpa_group_rekey=0 disables GTK rotation (each rekey is
    # a known reauth-bug trigger on old Android TV chipsets); uapsd=0
    # withdraws the WMM-PS advertisement (some STBs interpret it then sleep
    # too deep to receive the wake-up). All three are needed together;
    # adjusting just one didn't stop the disconnects in field testing.
    uci set wireless.${iface}.max_inactivity='31536000'
    uci set wireless.${iface}.skip_inactivity_poll='1'
    uci set wireless.${iface}.wpa_group_rekey='0'
    uci set wireless.${iface}.uapsd='0'
    # PMF (802.11w) optional — required by SAE (WPA3); WPA2 STAs ignore.
    # Old Android TV / IPTV STBs that crashed on PMF=2 are OK with PMF=1
    # (optional means it's negotiated only if STA advertises support).
    # If a specific STB still has issues, drop ieee80211w to 0 + use
    # encryption='psk2' (loses WPA3 but maximally compatible).
    uci set wireless.${iface}.ieee80211w='1'
    uci set wireless.${iface}.wmm='1'
    # 802.11k/v/r off — old STBs misinterpret BSS-Transition / FT requests
    # and sometimes treat them as forced disconnects.
    uci set wireless.${iface}.ieee80211r='0'
    uci set wireless.${iface}.ieee80211k='0'
    uci set wireless.${iface}.ieee80211v='0'
    uci -q delete wireless.${iface}.basic_rate
    uci add_list wireless.${iface}.basic_rate='6000'
    uci add_list wireless.${iface}.basic_rate='12000'
    uci add_list wireless.${iface}.basic_rate='24000'
    uci set wireless.${iface}.mcast_rate='24000'
done

# Commit and bring up WiFi
uci commit wireless
wifi up

# Reload network
/etc/init.d/network reload

exit 0
