# files/etc/init.d/vipin-sta-keepalive
#!/bin/sh /etc/rc.common
# Wakes stuck WiFi STAs by deauth — see /usr/sbin/vipin-sta-keepalive header.
START=99
USE_PROCD=1

start_service() {
    procd_open_instance vipin-sta-keepalive
    procd_set_param command /usr/sbin/vipin-sta-keepalive
    procd_set_param respawn 5 5 0
    procd_set_param stdout 0
    procd_set_param stderr 0
    procd_close_instance
}
