# files/etc/init.d/vipin-dns-reset
#!/bin/sh /etc/rc.common
# Reset dnsmasq DNS policy at every boot, BEFORE dnsmasq (S19) starts.
#
# Why: configure_dns() in vipin-vpn writes noresolv=1 + server=<CN DNS
# or dnscrypt> into UCI once the VPN tunnel is up. UCI is persistent,
# so a router that comes up without an immediately-usable VPN server
# (e.g. boot before DHCP, server went offline, or just normal reboot
# ordering) will have dnsmasq start with the last-session upstream
# locked in — which is dead until vipin-vpn re-establishes the tunnel.
#
# This creates the same chicken-and-egg as first-boot clock skew:
# dnsmasq can't answer → sysntpd can't resolve pool.ntp.org → NTP
# never syncs → (and the usual cascade). 44aeccfd fixed first-boot
# via zzz-defaults, but uci-defaults only runs once — reboots weren't
# covered.
#
# START=18 sits between S10boot and S19dnsmasq. We don't restart
# dnsmasq here because it hasn't started yet — just ensure UCI is in
# the right state when S19dnsmasq reads it.

START=18
STOP=

boot() {
    # Remove the two keys that turn dnsmasq into "dnscrypt-only" mode.
    # Absent keys → dnsmasq reads /tmp/resolv.conf.d/resolv.conf.auto
    # (ISP DNS from WAN DHCP) on top of whatever `server=` is set.
    uci -q delete dhcp.@dnsmasq[0].noresolv
    uci -q delete dhcp.@dnsmasq[0].resolvfile

    # Collapse any persisted upstream list back to the dnscrypt loopback.
    # dnscrypt-proxy boots independently (S15) and will answer once its
    # own clock/TLS are healthy. Meanwhile dnsmasq's WAN-DNS fallback
    # carries traffic.
    uci -q delete dhcp.@dnsmasq[0].server
    uci add_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'

    uci commit dhcp

    logger -t vipin-dns-reset "DNS policy reset to WAN-fallback + dnscrypt preferred"
}
