# files/etc/init.d/vipin-vpn
#!/bin/sh /etc/rc.common

START=90
STOP=10
USE_PROCD=1

NAME=vipin-vpn

# MOCK mode: functions can be sourced for tests without side effects.
VIPIN_VPN_MOCK="${VIPIN_VPN_MOCK:-0}"

resolve_connect_server() {
    # Returns the hostname stunnel should actually connect to.
    # UCI-selected server is what the user picked (and what we display);
    # if they picked cn.fanq.in but the router is outside CN, we connect
    # to un.fanq.in — an overseas frontend that forwards CN traffic to cn2.
    # Treats un as a cn alias so the rest of the codebase doesn't care.
    local router_country uci_server base_domain
    if [ "$VIPIN_VPN_MOCK" = "1" ]; then
        router_country="$VIPIN_VPN_ROUTER_COUNTRY"
        uci_server="$VIPIN_VPN_UCI_SERVER"
        base_domain="$VIPIN_VPN_BASE_DOMAIN"
    else
        # Use vipin-country-ips get-country to match the detection used in
        # start_service. vipin-detect with no subcommand prints usage, not JSON;
        # get-country returns just the 2-letter code (empty if not yet detected).
        router_country=$(/usr/sbin/vipin-country-ips get-country 2>/dev/null)
        [ -z "$router_country" ] && router_country="$(uci -q get vipin.vpn.country || echo cn)"
        uci_server=$(uci -q get vipin.vpn.server)
        base_domain=$(uci -q get vipin.vpn.base_domain || echo fanq.in)
    fi
    if [ "$router_country" != "cn" ] && \
       echo "$uci_server" | grep -q "^cn\.${base_domain}$"; then
        echo "un.${base_domain}"
    else
        echo "$uci_server"
    fi
}

fetch_auth_params() {
    # POST /api/v1/router-auth source=stunnel with creds, write response to
    # $VIPIN_CONFIG_DIR/auth-params.json.
    # Exit codes: 0=ok, 1=expired/invalid, 2=network error.
    local username="$1" password="$2" mac="$3"
    local cfgdir="${VIPIN_CONFIG_DIR:-/etc/vipin}"
    local api_base response status
    mkdir -p "$cfgdir"

    if [ "$VIPIN_VPN_MOCK" = "1" ]; then
        response="$VIPIN_VPN_API_RESPONSE"
    else
        api_base=$(uci -q get vipin.vpn.site_url || echo "https://www.anyfq.com")
        response=$(wget -q -O- --timeout=10 \
            "${api_base}/api/v1/router-auth" \
            --post-data="username=${username}&credentials=${password}&mac=${mac}&source=stunnel" \
            2>/dev/null)
    fi

    if [ -z "$response" ]; then
        return 2
    fi

    echo "$response" > "$cfgdir/auth-params.json"
    status=$(echo "$response" | grep -o '"status":"[^"]*"' | cut -d'"' -f4)
    case "$status" in
        ok) return 0 ;;
        expired|invalid|not_vip) return 1 ;;
        *) return 2 ;;
    esac
}

extract_ca_cert() {
    # Pull ca_cert field out of auth-params.json into a PEM file.
    local cfgdir="${VIPIN_CONFIG_DIR:-/etc/vipin}"
    local params="$cfgdir/auth-params.json"
    [ -f "$params" ] || return 1
    # The ca_cert value is JSON-string with \n escapes — unescape to real newlines.
    awk 'match($0, /"ca_cert":"[^"]*"/) {
        s = substr($0, RSTART+11, RLENGTH-12)
        gsub(/\\n/, "\n", s)
        print s
    }' "$params" > "$cfgdir/stunnel-ca.pem"
    [ -s "$cfgdir/stunnel-ca.pem" ]
}

render_stunnel_conf() {
    # Build the stunnel client config from auth-params.json + resolved server.
    local cfgdir="${VIPIN_CONFIG_DIR:-/etc/vipin}"
    local params="$cfgdir/auth-params.json"
    [ -f "$params" ] || return 1

    local port connect_host
    port=$(grep -o '"port":[0-9]*' "$params" | head -1 | cut -d: -f2)
    connect_host=$(resolve_connect_server)
    [ -z "$port" ] || [ -z "$connect_host" ] && return 1

    # Pull optional DNS sidecar endpoint from the server response. Used when
    # overseas-router + CN-server — see configure_dns decision matrix.
    local dns_host dns_port
    dns_host=$(grep -o '"stunnel_host":"[^"]*"' "$params" | cut -d'"' -f4)
    dns_port=$(grep -o '"stunnel_port":[0-9]*' "$params" | head -1 | cut -d: -f2)

    cat > "$cfgdir/stunnel-client.conf" <<EOF
; Generated by vipin-vpn start_service — do not edit by hand.
pid = /var/run/vipin-stunnel.pid
foreground = yes
; Aggressive memory discipline for 128 MB routers. Without these, stunnel
; grows past 30 MB RSS under LAN traffic and gets OOM-killed repeatedly.
debug = 2
output = /dev/null
; Per-thread stack 64 KB (default 2 MB) — biggest single saving.
stack = 65536
; Limit entropy collection to what OpenSSL actually needs.
RNDbytes = 256
socket = l:TCP_NODELAY=1
socket = r:TCP_NODELAY=1
socket = l:SO_KEEPALIVE=1
socket = r:SO_KEEPALIVE=1

[vipn]
client = yes
accept = 127.0.0.1:1080
connect = ${connect_host}:${port}
; Round-robin across every A-record DNS returned for connect_host. stunnel's
; default is 'prio' which sticks to the first resolved IP and only fails over
; on socket error — a load-balanced origin (jp/hk have 4-5 A records each)
; would see all its traffic hit one backend. Verified live on a router: 6
; parallel tunnels hit 1 IP with prio, 4 distinct IPs with rr. DNS is still
; cached at stunnel startup (delay=no default), so this requires DNS already
; ready when start_service runs — which is guaranteed by the time we get
; here since fetch_auth_params just completed a successful API round-trip.
failover = rr
CAfile = ${cfgdir}/stunnel-ca.pem
verifyChain = yes
checkHost = ${connect_host}
; keep idle tunnels short so the process footprint doesn't grow under burst
TIMEOUTclose = 0
TIMEOUTidle = 60
EOF

}

render_stubby_conf() {
    # Emit /etc/vipin/stubby.yml pointed at the server-delivered DoT endpoint
    # (dns.stunnel_host:dns.stunnel_port) with SPKI pinning against the CA
    # cert we already pulled for the main tunnel. stubby (DoT client) handles
    # the UDP↔TCP framing dnsmasq expects, unlike a raw socat passthrough.
    local cfgdir="${VIPIN_CONFIG_DIR:-/etc/vipin}"
    local params="$cfgdir/auth-params.json"
    [ -f "$params" ] || return 1

    local dns_host dns_port
    dns_host=$(grep -o '"stunnel_host":"[^"]*"' "$params" | cut -d'"' -f4)
    dns_port=$(grep -o '"stunnel_port":[0-9]*' "$params" | head -1 | cut -d: -f2)
    [ -z "$dns_host" ] || [ -z "$dns_port" ] && return 1

    [ -x /usr/bin/openssl ] || return 1
    [ -s "$cfgdir/stunnel-ca.pem" ] || return 1

    # SPKI pin (base64 SHA-256 of the DER-encoded SubjectPublicKeyInfo).
    # Works with self-signed certs where full chain validation would fail.
    local pin
    pin=$(openssl x509 -in "$cfgdir/stunnel-ca.pem" -noout -pubkey 2>/dev/null \
          | openssl pkey -pubin -outform der 2>/dev/null \
          | openssl dgst -sha256 -binary 2>/dev/null \
          | openssl base64 2>/dev/null)
    [ -z "$pin" ] && return 1

    # Resolve the host so stubby doesn't depend on dnsmasq (which we're
    # about to reconfigure). Emit EVERY A record as its own upstream so
    # round_robin_upstreams below actually has something to rotate across
    # when the DoT frontend sits behind a multi-IP DNS record. Today
    # cn.fanq.in is single-IP (this loop emits one upstream), but a future
    # CDN/multi-homing rollout works out of the box.
    local ips
    ips=$(nslookup "$dns_host" 2>/dev/null | awk '/^Address.*:/ && $NF ~ /^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$/ { print $NF }')
    [ -z "$ips" ] && return 1

    cat > "$cfgdir/stubby.yml" <<EOF
# Generated by vipin-vpn start_service — do not edit by hand.
resolution_type: GETDNS_RESOLUTION_STUB
dns_transport_list:
  - GETDNS_TRANSPORT_TLS
tls_authentication: GETDNS_AUTHENTICATION_REQUIRED
tls_query_padding_blocksize: 128
edns_client_subnet_private: 1
# 60s idle — stubby's default inherited from stock OpenWrt (9s) was tuned
# for public DoT providers (Cloudflare/Quad9) where short connections help
# load-balancing. We point at our own cn.fanq.in:994 with only one client
# (this router) so there's no server-side reason to churn. A bursty browse
# (page load → 20 queries, pause → click next) pays a full TLS reconnect
# every time the pause exceeds the timeout; bumping 10s → 60s cuts cold
# reconnects ~6× on typical traffic. Cloudflare's published guidance is 30s.
idle_timeout: 60000
listen_addresses:
  - 127.0.0.1@5356
round_robin_upstreams: 1
upstream_recursive_servers:
EOF
    for ip in $ips; do
        cat >> "$cfgdir/stubby.yml" <<EOF
  - address_data: ${ip}
    tls_port: ${dns_port}
    tls_auth_name: "${dns_host}"
    tls_pubkey_pinset:
      - digest: "sha256"
        value: ${pin}
EOF
    done
}

configure_dns() {
    # Decide which dnsmasq upstream to use and (in non-MOCK mode) write it.
    # Prints the chosen upstream server string to stdout either way.
    local router_country uci_server base_domain server_country
    if [ "$VIPIN_VPN_MOCK" = "1" ]; then
        router_country="$VIPIN_VPN_ROUTER_COUNTRY"
        uci_server="$VIPIN_VPN_UCI_SERVER"
        base_domain="$VIPIN_VPN_BASE_DOMAIN"
    else
        router_country=$(/usr/sbin/vipin-country-ips get-country 2>/dev/null)
        [ -z "$router_country" ] && router_country="$(uci -q get vipin.vpn.country || echo cn)"
        uci_server=$(uci -q get vipin.vpn.server)
        base_domain=$(uci -q get vipin.vpn.base_domain || echo fanq.in)
    fi

    # un is a cn alias for split-tunnel and DNS purposes.
    server_country=$(echo "$uci_server" | cut -d. -f1)
    [ "$server_country" = "un" ] && server_country="cn"

    # 4-matrix decision:
    #   overseas router + cn server → stubby DoT sidecar at 127.0.0.1:5356
    #     (→ cn.fanq.in:994 DoT → cn2 dnsmasq → 223.5.5.5). Returns CN-
    #     perspective answers so baidu resolves to mainland IPs.
    #   everything else → dnscrypt-proxy on 127.0.0.1:5353.
    # Sidecar active only when render_stubby_conf produced /etc/vipin/stubby.yml.
    local upstream="127.0.0.1#5353"
    if [ "$router_country" != "cn" ] && [ "$server_country" = "cn" ] && \
       [ -s "${VIPIN_CONFIG_DIR:-/etc/vipin}/stubby.yml" ]; then
        upstream="127.0.0.1#5356"
    fi

    if [ "$VIPIN_VPN_MOCK" != "1" ]; then
        uci -q delete dhcp.@dnsmasq[0].server
        uci add_list dhcp.@dnsmasq[0].server="$upstream"
        uci set dhcp.@dnsmasq[0].noresolv='1'
        uci commit dhcp
        /etc/init.d/dnsmasq restart >/dev/null 2>&1
    fi
    echo "$upstream"
}

teardown_dns() {
    # Revert to dnscrypt-only upstream when VPN is stopped.
    if [ "$VIPIN_VPN_MOCK" != "1" ]; then
        uci -q delete dhcp.@dnsmasq[0].server
        uci add_list dhcp.@dnsmasq[0].server='127.0.0.1#5353'
        uci set dhcp.@dnsmasq[0].noresolv='1'
        uci commit dhcp
        /etc/init.d/dnsmasq restart >/dev/null 2>&1
    fi
}

get_config() {
    config_load vipin
    config_get_bool split_tunnel vpn split_tunnel 1
    config_get server vpn server ""
    config_get username vpn username ""
    config_get country vpn country "cn"
    config_get base_domain vpn base_domain "fanq.in"
    config_get site_url vpn site_url "https://www.anyfq.com"
}

check_auth() {
    local status
    status=$(/usr/sbin/vipin-auth status 2>/dev/null)
    echo "$status" | grep -q '"logged_in": *true'
}

start_service() {
    get_config

    check_auth || return 1

    mkdir -p /var/log /etc/vipin /tmp/vipin

    local auth_info username pass_file
    auth_info=$(/usr/sbin/vipin-auth status 2>/dev/null)
    username=$(echo "$auth_info" | grep -o '"username": *"[^"]*"' | cut -d'"' -f4)
    pass_file="/etc/vipin/.vpn_pass"

    if [ -z "$username" ] || [ ! -f "$pass_file" ]; then
        logger -t "$NAME" "VPN not configured, skipping start"
        return 1
    fi

    local password mac
    password=$(cat "$pass_file")
    mac=$(/usr/sbin/vipin-auth mac 2>/dev/null | \
          grep -o '"mac": *"[^"]*"' | cut -d'"' -f4)

    # Ensure DNS is at a non-stunnel-dependent upstream before reaching out to
    # the auth API — configure_dns may have previously set the upstream to
    # 127.0.0.1#5356 (the stunnel bridge), which is dead whenever the VPN
    # stack isn't running. teardown_dns reverts to dnscrypt (127.0.0.1#5353),
    # which works standalone.
    teardown_dns

    fetch_auth_params "$username" "$password" "$mac"
    case $? in
        0) ;;  # ok — proceed
        1) logger -t "$NAME" "Auth expired/invalid, not starting VPN"
           uci set vipin.vpn.auth_status="$(grep -o '"status":"[^"]*"' /etc/vipin/auth-params.json | cut -d'"' -f4)"
           uci commit vipin
           return 1 ;;
        2) logger -t "$NAME" "Auth API unreachable on start, deferring"
           return 1 ;;
    esac

    uci set vipin.vpn.auth_status='ok'
    uci commit vipin

    extract_ca_cert || {
        logger -t "$NAME" "Failed to extract CA cert"
        return 1
    }
    render_stunnel_conf || {
        logger -t "$NAME" "Failed to render stunnel config"
        return 1
    }

    ip link delete tun0 2>/dev/null

    local connect_host
    connect_host=$(resolve_connect_server)
    logger -t "$NAME" "Starting stunnel client → ${connect_host}:$(grep -o '"port":[0-9]*' /etc/vipin/auth-params.json | head -1 | cut -d: -f2)"

    procd_open_instance stunnel
    procd_set_param command /usr/bin/stunnel /etc/vipin/stunnel-client.conf
    procd_set_param respawn 5 30 3
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance

    # DNS-over-TLS sidecar: stubby on 127.0.0.1:5356 → cn.fanq.in:994.
    # stubby handles UDP↔TCP framing + TLS + SPKI pinning natively. Only
    # launched when render_stubby_conf actually produced /etc/vipin/stubby.yml
    # (server returned dns.stunnel_host, cert present, openssl present).
    render_stubby_conf
    if [ -x /usr/sbin/stubby ] && [ -s /etc/vipin/stubby.yml ]; then
        procd_open_instance stubby
        procd_set_param command /usr/sbin/stubby -C /etc/vipin/stubby.yml
        procd_set_param respawn 5 30 3
        procd_close_instance
    fi

    cat > /etc/vipin/hev-tun.yaml << 'YAML'
tunnel:
  name: tun0
  mtu: 8500
  multi-queue: true
  ipv4: 192.168.200.1/24
socks5:
  address: 127.0.0.1
  port: 1080
  udp: 'tcp'
misc:
  task-stack-size: 20480
  log-file: stderr
  log-level: warn
YAML

    procd_open_instance tun2socks
    procd_set_param command /usr/bin/hev-socks5-tunnel /etc/vipin/hev-tun.yaml
    procd_set_param respawn 5 30 3
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance

    # Poll for tun0 — hev-socks5-tunnel brings the device up only after
    # stunnel has handshaken, which on MT7628-class boxes can take 20s+
    # under WAN latency. If the poll expires before tun0 shows up,
    # vipin-auth-check's 30s watchdog reinstates the default routes
    # as soon as the device appears.
    ( i=0
      while [ "$i" -lt 30 ]; do
          ip link show tun0 >/dev/null 2>&1 && break
          sleep 1
          i=$((i + 1))
      done
      if ip link show tun0 >/dev/null 2>&1; then
          ip link set tun0 up 2>/dev/null
          ip route add 0.0.0.0/1 dev tun0 2>/dev/null
          ip route add 128.0.0.0/1 dev tun0 2>/dev/null
          local split_mode country dns_upstream
          split_mode=$(uci -q get vipin.vpn.split_mode || echo forward)
          if [ "$split_mode" = "reverse" ]; then
              # reverse: server-country IPs go through VPN (so user picks up
              # that country's egress). un is a cn alias.
              country=$(echo "$(uci -q get vipin.vpn.server)" | cut -d. -f1)
              [ "$country" = "un" ] && country="cn"
              /usr/sbin/vipin-vpn-routing enable "$country" reverse
          elif [ "$split_mode" = "forward" ]; then
              # forward: router-location IPs stay on direct WAN, the rest
              # goes through VPN. Use the auto-detected router country, not
              # the static uci vpn.country (which is a server-side hint).
              country=$(/usr/sbin/vipin-country-ips get-country 2>/dev/null)
              [ -z "$country" ] && country="$(uci -q get vipin.vpn.country || echo cn)"
              /usr/sbin/vipin-vpn-routing enable "$country" forward
          fi
          # stack is up — now safe to flip DNS to the 4-matrix decision
          # (stunnel-dependent bridge only becomes usable at this point).
          dns_upstream=$(configure_dns)
          logger -t "$NAME" "DNS upstream set to $dns_upstream"
          logger -t "$NAME" "VPN started"
      else
          logger -t "$NAME" "tun0 did not come up in 30s"
      fi
    ) &
}

stop_service() {
    get_config

    /usr/sbin/vipin-vpn-routing disable 2>/dev/null

    killall stunnel 2>/dev/null
    killall hev-socks5-tunnel 2>/dev/null
    sleep 1

    ip route del 0.0.0.0/1 dev tun0 2>/dev/null
    ip route del 128.0.0.0/1 dev tun0 2>/dev/null
    ip link delete tun0 2>/dev/null

    teardown_dns

    rm -f /var/run/vipin-stunnel.pid

    logger -t "$NAME" "VPN stopped"
}

_restart_locked() {
    stop_service
    sleep 1
    start_service
}

reload_service() {
    # Non-blocking lock: if another restart is already in progress, drop this
    # one instead of queueing. Concurrent restart storms (LuCI + watchdog +
    # user) would otherwise pile up on procd_lock for minutes.
    ( flock -n 8 || { logger -t "$NAME" "restart already in progress, skipping"; exit 0; }
      _restart_locked
    ) 8>/var/lock/vipin-restart.lock
}

restart() {
    ( flock -n 8 || { logger -t "$NAME" "restart already in progress, skipping"; exit 0; }
      _restart_locked
    ) 8>/var/lock/vipin-restart.lock
}

status() {
    # With USE_PROCD=1, rc.common overrides this function with procd's own
    # wrapper (which prints just "running"). Keep this for the direct-invoke
    # case and for shells without procd — it's never the LuCI code path.
    if pgrep stunnel >/dev/null && pgrep hev-socks5-tunnel >/dev/null && \
       ip link show tun0 >/dev/null 2>&1; then
        echo "VPN is running"
        local auth_info username
        auth_info=$(/usr/sbin/vipin-auth status 2>/dev/null)
        username=$(echo "$auth_info" | grep -o '"username": *"[^"]*"' | cut -d'"' -f4)
        [ -n "$username" ] && echo "User: $username"
        return 0
    fi
    echo "VPN is not running"
    return 1
}

boot() {
    # Start in background 30s after boot; gated by presence of auth files.
    ( sleep 30
      if [ -f /etc/vipin/auth.conf ] && [ -f /etc/vipin/.vpn_pass ]; then
          /etc/init.d/vipin-auth start 2>/dev/null
          start_service
      fi
    ) &
}
