# files/etc/init.d/vipin-video
#!/bin/sh /etc/rc.common
# vipin-video — video_direct feature init.d (runs enable on boot)

START=91
USE_PROCD=1
NAME=vipin-video

start_service() {
    /usr/sbin/vipin-video-domains enable 2>/dev/null || true
}

stop_service() {
    /usr/sbin/vipin-video-domains disable 2>/dev/null || true
}

reload_service() {
    /usr/sbin/vipin-video-domains enable 2>/dev/null || true
}
