# files/etc/init.d/vipin-vpn
#!/bin/sh /etc/rc.common

START=90
STOP=10
USE_PROCD=1

NAME=vipin-vpn
PROG=/usr/sbin/openconnect

# MOCK mode: functions can be sourced for tests without side effects.
VIPIN_VPN_MOCK="${VIPIN_VPN_MOCK:-0}"

resolve_connect_server() {
    # Returns the hostname stunnel should actually connect to.
    # UCI-selected server is what the user picked (and what we display);
    # if they picked cn.fanq.in but the router is outside CN, we connect
    # to un.fanq.in — an overseas frontend that forwards CN traffic to cn2.
    # Treats un as a cn alias so the rest of the codebase doesn't care.
    local router_country uci_server base_domain
    if [ "$VIPIN_VPN_MOCK" = "1" ]; then
        router_country="$VIPIN_VPN_ROUTER_COUNTRY"
        uci_server="$VIPIN_VPN_UCI_SERVER"
        base_domain="$VIPIN_VPN_BASE_DOMAIN"
    else
        # Use vipin-country-ips get-country to match the detection used in
        # start_service. vipin-detect with no subcommand prints usage, not JSON;
        # get-country returns just the 2-letter code (empty if not yet detected).
        router_country=$(/usr/sbin/vipin-country-ips get-country 2>/dev/null)
        [ -z "$router_country" ] && router_country="$(uci -q get vipin.vpn.country || echo cn)"
        uci_server=$(uci -q get vipin.vpn.server)
        base_domain=$(uci -q get vipin.vpn.base_domain || echo fanq.in)
    fi
    if [ "$router_country" != "cn" ] && \
       echo "$uci_server" | grep -q "^cn\.${base_domain}$"; then
        echo "un.${base_domain}"
    else
        echo "$uci_server"
    fi
}

get_config() {
    config_load vipin
    config_get_bool enabled vpn enabled 0
    config_get_bool split_tunnel vpn split_tunnel 1
    config_get server vpn server ""
    config_get username vpn username ""
    config_get country vpn country "cn"
    config_get base_domain vpn base_domain "fanq.in"
    config_get site_url vpn site_url "https://www.anyfq.com"
}

check_auth() {
    local status=$(/usr/sbin/vipin-auth status 2>/dev/null)
    
    if echo "$status" | grep -q '"logged_in": *true'; then
        return 0
    fi
    
    logger -t $NAME "Authentication required, cannot start VPN"
    return 1
}

probe_cert_pin() {
    local srv="$1"
    local pin_file="/etc/vipin/cert_pin"

    # Intentionally pass wrong pin so openconnect prints the real one
    local err=$(echo "" | openconnect --protocol=anyconnect --server "$srv:443" \
        --servercert pin-sha256:AAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAA= \
        --passwd-on-stdin 2>&1 | grep "pin-sha256:")
    local real_pin=$(echo "$err" | grep -o 'pin-sha256:[A-Za-z0-9+/=]*')
    if [ -n "$real_pin" ]; then
        echo "$real_pin" > "$pin_file"
        logger -t $NAME "Probed real cert pin for $srv: $real_pin"
        return 0
    fi
    return 1
}

fetch_cert_pin() {
    local srv="$1"
    local force="$2"
    local pin_file="/etc/vipin/cert_pin"

    # force=1: skip API, probe real pin directly (self-heal path)
    if [ "$force" = "1" ]; then
        probe_cert_pin "$srv"
        return $?
    fi

    # Try API first
    local response=$(wget -q -O- --timeout=10 "$SERVERS_API" 2>/dev/null)
    if [ -n "$response" ]; then
        local pin=$(echo "$response" | grep -o "\"${srv}\"[^}]*\"cert_pin\":\"[^\"]*\"" | grep -o '"cert_pin":"[^"]*"' | cut -d'"' -f4)
        if [ -n "$pin" ]; then
            echo "$pin" > "$pin_file"
            return 0
        fi
    fi

    # API unreachable or server missing: probe only if we have no pin at all
    [ ! -f "$pin_file" ] && probe_cert_pin "$srv"
    return $?
}

get_cert_pin() {
    local pin_file="/etc/vipin/cert_pin"
    if [ -f "$pin_file" ]; then
        cat "$pin_file"
    else
        echo ""
    fi
}

start_service() {
    get_config
    
    [ "$enabled" = "0" ] && return 1
    
    if ! check_auth; then
        return 1
    fi
    
    mkdir -p /var/log
    mkdir -p /etc/vipin
    mkdir -p /tmp/vipin
    
    local auth_info=$(/usr/sbin/vipin-auth status 2>/dev/null)
    local username=$(echo "$auth_info" | grep -o '"username": *"[^"]*"' | cut -d'"' -f4)
    local pass_file="/etc/vipin/.vpn_pass"

    if [ -z "$server" ] || [ -z "$username" ]; then
        logger -t $NAME "VPN not configured, skipping start"
        return 1
    fi

    if [ ! -f "$pass_file" ]; then
        logger -t $NAME "VPN password file not found, re-login required"
        return 1
    fi

    logger -t $NAME "Starting VPN for user: $username"
    
    # Detect router's current country
    local current_country=$(/usr/sbin/vipin-country-ips get-country 2>/dev/null)
    [ -n "$current_country" ] && country="$current_country"
    
    # If router is not in China and trying to connect to cn server,
    # use un.fanq.in for connection only (keeps cn.fanq.in in UCI for display/split)
    local connect_server="$server"
    if [ "$country" != "cn" ] && echo "$server" | grep -q "cn\.${base_domain}"; then
        logger -t $NAME "Router not in China ($country), connecting via un server (display stays cn)"
        connect_server="un.${base_domain}"
    fi
    
    /usr/sbin/vipin-country-ips update "$country" 2>/dev/null
    
    if pgrep openconnect >/dev/null; then
        logger -t $NAME "VPN already running"
        return 0
    fi

    # Clean up stale tun device to prevent TUNSETIFF: Resource busy
    ip link delete vpn_vipin 2>/dev/null
    
    local vpn_port=$((50000 + RANDOM % 10000))

    fetch_cert_pin "$connect_server"
    local cert_pin=$(get_cert_pin)
    local cert_opt=""
    if [ -n "$cert_pin" ]; then
        cert_opt="--servercert $cert_pin"
    fi

    logger -t $NAME "Connecting to $connect_server:$vpn_port (config: $server)"

    # Create wrapper script to pipe password to openconnect via stdin
    # Note: the router's busybox-linked openconnect build does not expose
    # --dpd on the command line; DTLS keepalive is set entirely by the
    # server's ocserv.conf dpd value, which we tuned separately.
    cat > /tmp/vipin/vpn_start.sh << SCRIPT
#!/bin/sh
cat "$pass_file" | $PROG \
    --protocol=anyconnect \
    --server "$connect_server:$vpn_port" \
    --useragent "OpenConnect VPN Agent" \
    --authgroup "" \
    --user "$username" \
    --passwd-on-stdin \
    --script /usr/sbin/vipin-vpnc-script \
    $cert_opt \
    --interface vpn_vipin \
    --pid-file /var/run/vipin-vpn.pid \
    2>/var/log/openconnect.log
SCRIPT
    chmod +x /tmp/vipin/vpn_start.sh

    procd_open_instance
    procd_set_param command /bin/sh /tmp/vipin/vpn_start.sh
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance

    # Check connection status in background; self-heal on pin mismatch
    ( sleep 10
      if pgrep openconnect >/dev/null; then
          logger -t $NAME "VPN started (server: $server)"
          rm -f /tmp/vipin/.pin_retried
      else
          logger -t $NAME "Failed to start VPN"
          tail -5 /var/log/openconnect.log 2>/dev/null | logger -t $NAME

          # Detect pin mismatch and self-heal (only retry once per start cycle)
          if grep -q -E "fingerprint.*match|signer not found|certificate verify failed" /var/log/openconnect.log 2>/dev/null \
             && [ ! -f /tmp/vipin/.pin_retried ]; then
              logger -t $NAME "Cert pin mismatch detected, probing real pin"
              touch /tmp/vipin/.pin_retried
              rm -f /etc/vipin/cert_pin
              if fetch_cert_pin "$server" 1; then
                  logger -t $NAME "Real pin acquired, restarting VPN"
                  /etc/init.d/vipin-vpn restart
              else
                  logger -t $NAME "Failed to probe real pin, giving up"
              fi
          fi
      fi
    ) &
}

stop_service() {
    get_config
    
    /usr/sbin/vipin-vpn-routing disable 2>/dev/null
    
    if pgrep openconnect >/dev/null; then
        killall openconnect 2>/dev/null
        sleep 2
    fi
    
    rm -f /var/run/vipin-vpn.pid
    
    logger -t $NAME "VPN stopped"
}

reload_service() {
    stop_service
    sleep 1
    start_service
}

restart() {
    stop_service
    sleep 1
    start_service
}

status() {
    if pgrep openconnect >/dev/null; then
        echo "VPN is running"
        local auth_info=$(/usr/sbin/vipin-auth status 2>/dev/null)
        local username=$(echo "$auth_info" | grep -o '"username": *"[^"]*"' | cut -d'"' -f4)
        if [ -n "$username" ]; then
            echo "User: $username"
        fi
        return 0
    else
        echo "VPN is not running"
        return 1
    fi
}

boot() {
    get_config

    [ "$enabled" = "1" ] && {
        # Wait for network in background to avoid blocking other init scripts
        ( sleep 30
          /etc/init.d/vipin-auth start 2>/dev/null
          start_service
        ) &
    }
}
