# files/etc/hotplug.d/iface/50-vipin-bypass
#!/bin/sh
# Install /32 bypass routes for dnscrypt-proxy's DoH upstream IPs so they
# ALWAYS resolve via WAN, never via tun0. Without these, vipin-vpn's
# 0.0.0.0/1 dev tun0 catchall swallows dnscrypt's TLS handshakes to
# Cloudflare/Google/Quad9:443 — and if tun0 is zombie (stunnel handshake
# never completed) dnscrypt can never bootstrap → VPN can never rely on
# dnscrypt for CN-DNS either → whole DNS stack dies.
#
# These IPs match dnscrypt-proxy.toml's server_names (cloudflare / google)
# sdns:// stamps and fallback_resolvers / netprobe_address. Keep in sync
# with files/etc/dnscrypt-proxy/dnscrypt-proxy.toml.

[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] || exit 0

# Grab the just-provisioned WAN gateway + device.
wan_gw=$(ip route show default 2>/dev/null | awk '/default/ {print $3; exit}')
wan_dev=$(ip route show default 2>/dev/null | awk '/default/ {print $5; exit}')
[ -z "$wan_gw" ] || [ -z "$wan_dev" ] && exit 0

# DoH upstreams dnscrypt tries to reach at startup + netprobe + fallback.
# Using `ip route replace` so re-fires on WAN flap update cleanly rather
# than failing-loud with "file exists".
for ip in 1.1.1.1 1.0.0.1 8.8.8.8 8.8.4.4 9.9.9.9 149.112.112.112; do
    ip route replace "$ip/32" via "$wan_gw" dev "$wan_dev" 2>/dev/null
done

logger -t vipin-bypass "dnscrypt DoH bypass routes installed via $wan_gw dev $wan_dev"
