# files/etc/hotplug.d/iface/89-vipin-ntp-bootstrap
#!/bin/sh
# First-boot NTP bootstrap via hardcoded anycast IPs (no DNS required).
#
# Problem this solves: routers without an RTC battery power up with a
# kernel-build-date clock (e.g. 2026-03-17). TLS certs issued later than
# that date fail x509 NotBefore validation, which breaks dnscrypt-proxy
# (DoH to 8.8.8.8 / 1.1.1.1). No dnscrypt → no DNS → stock sysntpd can't
# resolve pool.ntp.org → clock stays broken → chicken-and-egg deadlock.
#
# Fix: on WAN ifup, if system year < 2025 (obviously wrong), force an
# SNTP sync against known anycast NTP IPs, bypassing DNS entirely. Must
# run before 99-dnscrypt-ready so dnscrypt retries with a valid clock.

[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] || exit 0

FLAG="/tmp/.ntp-bootstrap-done"
[ -f "$FLAG" ] && exit 0

year=$(date +%Y)
if [ "$year" -ge 2025 ] 2>/dev/null; then
    touch "$FLAG"
    exit 0
fi

log() { logger -t vipin-ntp-bootstrap "$*"; }

log "clock shows year=$year, forcing SNTP via anycast IPs"

# Anycast / globally-reachable NTP IPs. Kept small to keep ntpd -q fast.
NTP_IPS="162.159.200.1 216.239.35.0 203.107.6.88"

for ip in $NTP_IPS; do
    if ntpd -q -n -p "$ip" >/dev/null 2>&1; then
        log "time set via $ip → $(date -u '+%F %T UTC')"
        touch "$FLAG"
        # Kick dnscrypt so it retries TLS with the now-valid clock.
        # 99-dnscrypt-ready will run after us and handle dnsmasq reload.
        /etc/init.d/dnscrypt-proxy restart 2>/dev/null
        exit 0
    fi
done

log "all NTP IPs failed, clock still broken"
exit 1
