# files/etc/hotplug.d/iface/98-vipin-video-firstboot
#!/bin/sh
# One-shot: refresh the video_direct domain list from GitHub the first
# time WAN comes up after a flash/boot. Firmware ships with a baked-in
# copy, but the authoritative list evolves; without this hook users
# would keep the image-age list until they manually hit refresh.
#
# Gated on /tmp flag so it fires once per boot (watchdog restarts of
# WAN mid-session don't re-hammer GitHub). DNS can still be racing at
# ifup, so we wait for it to resolve before handing off.

[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] || exit 0

FLAG="/tmp/.vipin-video-refreshed"
[ -f "$FLAG" ] && exit 0
touch "$FLAG"

(
    # Wait up to 90s for DNS to resolve the GitHub host. Same failure
    # mode as 99-dnscrypt-ready: on cold flash, dnscrypt-proxy takes
    # ~60s to bootstrap, so nslookup returns NXDOMAIN until then.
    i=0
    while [ $i -lt 90 ]; do
        if nslookup raw.githubusercontent.com 127.0.0.1 >/dev/null 2>&1; then
            break
        fi
        sleep 3
        i=$((i + 3))
    done

    logger -t vipin-video-firstboot "DNS ready, refreshing video-domain list"
    if /usr/sbin/vipin-video-domains refresh >/dev/null 2>&1; then
        logger -t vipin-video-firstboot "refresh ok"
    else
        logger -t vipin-video-firstboot "refresh failed — will retry next boot"
        rm -f "$FLAG"
    fi
) &
