# files/etc/hotplug.d/iface/99-dnscrypt-ready
#!/bin/sh
# Wait for dnscrypt-proxy to be ready after WAN comes up, then restart dnsmasq
# Fixes: dnsmasq starts before dnscrypt-proxy is listening on 127.0.0.1:5353

[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] || exit 0

# Only run once per boot
FLAG="/tmp/.dnscrypt-ready-done"
[ -f "$FLAG" ] && exit 0

logger -t dnscrypt-ready "WAN up, waiting for dnscrypt-proxy on 127.0.0.1:5353..."

# Wait up to 60 seconds for dnscrypt-proxy to listen
i=0
while [ $i -lt 60 ]; do
    if netstat -tlnup 2>/dev/null | grep -q '127.0.0.1:5353'; then
        logger -t dnscrypt-ready "dnscrypt-proxy ready, restarting dnsmasq"
        /etc/init.d/dnsmasq restart
        touch "$FLAG"
        exit 0
    fi
    sleep 1
    i=$((i + 1))
done

logger -t dnscrypt-ready "Timeout waiting for dnscrypt-proxy"
